# docker/encrypt/supervisord.conf
[supervisord]
nodaemon=true

; Apache server
[program:apache2]
command=apache2-foreground
